# instance/crontab
*/5 * * * * /usr/bin/python3 /home/ubuntu/autoshutdown.py